# launchd/com.tfg.scraper.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tfg.scraper.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>/Users/david/Documents/scraper/scripts/run_daily_pipeline.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/david/Documents/scraper</string>

    <key>StandardOutPath</key>
    <string>/Users/david/Documents/scraper/data/raw/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/david/Documents/scraper/data/raw/launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
